# service/ec2-associate-address
#! /bin/sh
# chkconfig: 35 55 25
# description: ec2-associate-address

### config
eip=
###

. /etc/init.d/functions
. /etc/profile.d/aws-apitools-common.sh

if [ ! -f /etc/sysconfig/aws.env ]; then
    exit 6;
fi

. /etc/sysconfig/aws.env

start() {
    ec2-associate-address -i $instance_id $eip --region ap-northeast-1
}

stop() {
    ec2-disassociate-address $eip --region ap-northeast-1
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    sleep 3
    ;;
  restart)
    $0 stop
    $0 start
    ;;
  *)
    echo "Usage $0 [start|stop|restart]"
esac